# scripts/v2ray/sspanel_V2ray.service
[Unit]
Description=restart sspanel_V2ray service

[Service]
ExecStart=/bin/bash /root/v2scar/v2ray.sh restart

[Install]
WantedBy=multi-user.target

